# launchd/com.plo.autostockbot.startall.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.plo.autostockbot.startall</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/plo/Documents/auto_stock_bot/scripts/autostart_all.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/plo/Documents/auto_stock_bot</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/plo/Documents/auto_stock_bot/logs/launchd_autostart.out</string>

  <key>StandardErrorPath</key>
  <string>/Users/plo/Documents/auto_stock_bot/logs/launchd_autostart.err</string>
</dict>
</plist>
